# Systemd/system/syncoid.timer
[Unit]
Description=Starting syncoid daily

[Timer]
OnCalendar=Daily
AccuracySec=1h
Persistent=true

[Install]
WantedBy=multi-user.target
